# web/xsisi_pyramids/supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=python3 /app/app.py


[program:nodejs]
command=node /app/bot/index.js

